# systemd/user/stardeck-core.target
[Unit]
Description=Stardeck Core Services
Requires=dbus.service
Wants=plusdeck.service
Wants=crystalfontz.service
