# systemd/user/stardeck.target
[Unit]
Description=Stardeck Service
Requires=stardeck-core.target
Requires=stardeck.service
